# systemd/animus-hermes-healthcheck.service
[Unit]
Description=ANIMUS + Hermes hourly health check
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/animus
ExecStart=%h/animus/scripts/check-animus-hermes.sh
EnvironmentFile=-%h/animus/animus.env
