# systemd/animus-hermes-healthcheck.timer
[Unit]
Description=Run ANIMUS + Hermes health check hourly

[Timer]
OnBootSec=5min
OnCalendar=hourly
RandomizedDelaySec=2min
Unit=animus-hermes-healthcheck.service
Persistent=true

[Install]
WantedBy=timers.target
