# systemd/animus.service
[Unit]
Description=ANIMUS chat UI (restart: systemctl --user restart animus.service)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/animus/animus-chat
ExecStart=%h/animus/animus-chat/.venv/bin/python server.py
EnvironmentFile=-%h/animus/animus.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
